# templates/default/unicorn.upstart.conf.erb
#!upstart

env USER=<%= @user %>
env RAILS_ENV=<%= @environment %>

start on startup
stop  on shutdown

respawn

script
  chdir <%= @directory %>
  exec sudo -u <%= @user %> bundle exec unicorn -E <%= @environment %> -c <%= @config %>
end script

